# infra/com.llm-relay.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.llm-relay</string>

    <key>ProgramArguments</key>
    <array>
        <string>/opt/llm-relay/infra/start.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/opt/llm-relay</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/var/log/llm-relay/llm-relay.log</string>

    <key>StandardErrorPath</key>
    <string>/var/log/llm-relay/llm-relay-error.log</string>
</dict>
</plist>
