# infra/llm-relay.service
[Unit]
Description=LLM Relay Server
After=network.target

[Service]
Type=simple
User=llm-relay
WorkingDirectory=/opt/llm-relay
ExecStart=/usr/bin/node dist/index.js
EnvironmentFile=/opt/llm-relay/.env
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=llm-relay

[Install]
WantedBy=multi-user.target
